# examples/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.notion.auto-learn</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/notion-auto-learn/bin/process.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>/tmp/notion-auto-learn.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/notion-auto-learn-error.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
